# manifest/etc/cron.d/root
@hourly /usr/bin/dmarcts-report-parser.pl -i -d -r
