# manifest/etc/cron.d/root
15 * * * * /usr/bin/dmarcts-report-parser.pl -i -d -r 1>>/var/log/nginx/dmarc-reports.log 2>&1

