# manifest/etc/supervisor/conf.d/crond.conf
[program:cron]
command=crond -f -l 0 -L - -c /etc/cron.d
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
